# fpm-alpine/crontabs
* * * * * /var/cron/cron-entrypoint.sh >> /dev/null 2>&1
